# gitea_repos_master/isphere-services/proxy-check/config/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:proxy_check]
directory=/opt/proxy-check/src
command=sh -exc 'python proxy_check.py && sleep 60'
autorestart=false
autostart=true
exitcodes=0
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
user=root

[eventlistener:processes]
command=bash -c "printf 'READY\n' && while read line; do kill -SIGQUIT $PPID; done < /dev/stdin"
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL